# alarm-ioc/supervisord.conf
[supervisord]
nodaemon=true

[program:soft_ioc]
command=bash /opt/nalms/run.sh
