# examples/heart_of_gold/supervisord.conf

[supervisord]
nodaemon=true

[program:soft_ioc]
command=/bin/bash -c "softIoc -d /var/ioc/heart_of_gold.db"
